# dist/infowindow.timer
[Unit]
Description=InfoWindow updater

[Timer]
Unit=infowindow.service
OnCalendar=*:0/6

[Install]
WantedBy=timers.target
